# ext/ag_overlay_cleanup_cron.timer
[Unit]
Description=Timer for AG Overlay Cleanup Cron

[Timer]
OnUnitActiveSec=5s
AccuracySec=100ms
Persistent=true

[Install]
WantedBy=timers.target